# scripts/ztlab-posture.timer
[Unit]
Description=Run ZTLab posture check every 60 seconds

[Timer]
OnBootSec=30
OnUnitActiveSec=60

[Install]
WantedBy=timers.target
